# templates/backup.timer
[Unit]
Description=Run borgmatic backup for {{ APP_NAME }}

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target